# crontab
@reboot tmux new -d -s jupyter 'jupyter-lab --no-browser;'
